# templates/consul-template.sysv.erb
#!/bin/bash
#
#       /etc/rc.d/init.d/consul-template
#
#       Daemonize the consul-template agent.
#
# chkconfig:   2345 95 95
# description: Populate config files with the values in Consul's KV store
# processname: consul-template
# pidfile: /var/run/consul-template/pidfile

# Source function library.
. /etc/init.d/functions

CONSUL=<%= scope.lookupvar('consul_template::bin_dir') %>/consul-template
CONFIG=<%= scope.lookupvar('consul_template::config_dir') %>
PID_FILE=/var/run/consul-template/consul-template.pid
LOG_FILE=/var/log/consul-template

[ -e /etc/sysconfig/consul-template ] && . /etc/sysconfig/consul-template

mkrundir() {
        [ ! -d /var/run/consul-template ] && mkdir -p /var/run/consul-template
        chown <%= scope.lookupvar('consul_template::user') %> /var/run/consul-template
}

KILLPROC_OPT="-p ${PID_FILE}"
mkpidfile() {
        # Create PID file if it didn't exist
        mkrundir
        [ ! -f $PID_FILE ] && pidofproc $CONSUL > $PID_FILE
        chown <%= scope.lookupvar('consul_template::user') %> /var/run/consul-template
        if [ $? -ne 0 ] ; then
            rm $PID_FILE
            KILLPROC_OPT=""
        fi
}

start() {
        echo -n "Starting consul-template: "
        mkrundir
        [ -f $PID_FILE ] && rm $PID_FILE
        daemon --user=<%= scope.lookupvar('consul_template::user') %> \
            --pidfile="$PID_FILE" \
            "$CONSUL" -config "$CONFIG"/config.json <%= scope.lookupvar('consul_template::extra_options') %> >> "$LOG_FILE" &
        retcode=$?
        touch /var/lock/subsys/consul-template
        return $retcode
}

stop() {
        echo -n "Shutting down consul-template: "
        mkpidfile
        killproc $KILLPROC_OPT $CONSUL -9

        retcode=$?
        rm -f /var/lock/subsys/consul-template $PID_FILE
        return $retcode
}

status() {
    if [ -f "$PID_FILE" ]
    then
        if ps -p `cat $PID_FILE` > /dev/null 2&1
        then
            echo -n "Process running."
            return 0
        else
            echo -n "PID exists but process not running."
            return 1
        fi
    else
        echo -n "Process not running."
        return 3
    fi
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    status)
        status
        ;;
    restart)
        stop
        start
        ;;
    reload)
        kill -HUP `cat $PID_FILE`
        ;;
    condrestart)
        [ -f /var/lock/subsys/consul-template ] && restart || :
        ;;
    *)
        echo "Usage: consul-template {start|stop|status|reload|restart}"
        exit 1
        ;;
esac
exit $?
